# archcust/airootfs/etc/skel/.config/autostart/pasystray.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=pasystray
Icon=unknown--usr-bin-pasystray
TryExec=/usr/bin/pasystray
Exec=/usr/bin/pasystray
