# archcust/airootfs/etc/skel/.config/autostart/qxkb.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=qxkb
Icon=unknown--usr-bin-qxkb
TryExec=/usr/bin/qxkb
Exec=/usr/bin/qxkb
